# startup/ipcam.ini
[program:ipcam]
command=python3.2 /path/to/ipcam.py <token> <user_id>
;stdout_logfile=/path/to/stdout/log/file
;stderr_logfile=/path/to/stderr/log/file
